# supervisor.conf
[program:short]
command=uvicorn main:app
directory=/home/short/app/
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true